# certbot.cron
#!/bin/sh

certbot-renew -q
haproxy-refresh
